# privatecontributions.jmorais.github.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>privatecontributions.jmorais.github</string>
        <key>ProgramArguments</key>
        <array>
            <string>PROGRAM_PATH</string>
        </array>
        <key>StartInterval</key>
        <integer>43200</integer>
        <key>RunAtLoad</key>
        <false/>
    </dict>
</plist>